# motd/mira/update-motd.d/99-chafa
#!/bin/bash

COLORTERM=truecolor /usr/bin/chafa /opt/motd.png
