# motd/mira/update-motd.d/99-viu
#!/bin/bash

COLORTERM=truecolor viu /opt/motd.png
