# motd/aoi/update-motd.d/99-chafa
#!/bin/sh

/usr/bin/chafa /opt/aoi-s3.png
