# motd/aoi/update-motd.d/99-viu
#!/bin/sh

/usr/local/bin/viu /opt/aoi-s3.png
