# dotfile/etc/systemd/system/adguardhome.service.d/override.conf
[Service]
Restart=always
KillSignal=SIGKILL
